# claudecluster.service
[Unit]
Description=Claude Cluster - Distributed compute mesh for Claude sessions
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=paschal
WorkingDirectory=/home/paschal/claudecluster
ExecStart=/usr/bin/node /home/paschal/claudecluster/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
